# PI/systemd/controlapp.service
[Unit]
Description=Control App Main Service
After=network.target pigpiod.service

[Service]
User=controlstation
WorkingDirectory=/home/controlstation/app/PI
ExecStart=/home/controlstation/app/PI/.venv/bin/python /home/controlstation/app/PI/main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
