# PI/systemd/update-repo.timer
[Unit]
Description=Run repo update every 10 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=10min
Persistent=true

[Install]
WantedBy=timers.target
